# truevision-server.service
[Unit]
Description=TrueVision Server
After=network.target ollama.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/TrueVision
ExecStart=/home/ubuntu/TrueVision/venv/bin/uvicorn server.app:app --host 0.0.0.0 --port 8008
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=TRUEVISION_SERVER_PORT=8008

[Install]
WantedBy=multi-user.target
